# update_checker.service
[Unit]
Description=GitHub Update Checker
After=network.target

[Service]
User=www-data
WorkingDirectory=/var/www/html/CornGetFromGit
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 check_updates.py
Restart=always

[Install]
WantedBy=multi-user.target 