# update_checker.service
[Unit]
Description=GitHub Update Checker
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/update_checker
Environment="PATH=/usr/local/bin"
ExecStart=/usr/local/bin/python3 check_updates.py
Restart=always

[Install]
WantedBy=multi-user.target 